# scripts/com.agentregistry.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Service identifier -->
    <key>Label</key>
    <string>com.agentregistry.daemon</string>
    
    <!-- Command: bash wrapper to set ulimit before launching bun -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>ulimit -n 65536; exec "__BUN_PATH__" run "__AGENTREGISTRY_PATH__/src/server.ts" --daemon</string>
    </array>
    
    <!-- Working directory (use ~/.agentregistry, NOT project dir — launchd can't access ~/Documents) -->
    <key>WorkingDirectory</key>
    <string>__HOME__/.agentregistry</string>
    
    <!-- File descriptor limits (Bun requires more than macOS default 256) -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>
    <key>HardResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
    </dict>
    
    <!-- Start at login -->
    <key>RunAtLoad</key>
    <true/>
    
    <!-- Keep alive - restart on crash -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <!-- Throttle restarts (wait 10 seconds between restarts) -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>__HOME__/.agentregistry/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.agentregistry/logs/launchd-stderr.log</string>
    
    <!-- Environment -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:__HOME__/.bun/bin</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>AGENTREGISTRY_HOME</key>
        <string>__HOME__/.agentregistry</string>
    </dict>
</dict>
</plist>
